# news-scraper/crontab
*/2 * * * * python /app/news-scraper/main.py >> logs
